# supervisord/memcached.conf
[program:memcached]
command=/start-memcached.sh
numprocs=1
autostart=true
autorestart=true